# logitech_k810/50-k810.rules
KERNEL=="hidraw*", SUBSYSTEM=="hidraw", ATTRS{address}=="00:1f:20:a7:d5:01" \
 RUN+="/opt/logitech-k810/k810-udev.sh %p"
