# logitech_k810/50-k810.rules
KERNEL=="hidraw*", SUBSYSTEMS=="hidraw", RUN+="/opt/logitech-k810/k810-udev.sh %p"
