# logitech_k380/50-k380.rules
KERNEL=="hidraw*", SUBSYSTEMS=="hidraw", RUN+="/opt/logitech-k380/k380-udev.sh %p"
